# setup/debian/DEBIAN/postrm
#!/bin/bash
set -e

LOGFILE="/var/log/threepics-postrm.log"
exec >> "$LOGFILE" 2>&1

echo "➡ postrm started at $(date) with argument: $1"

case "$1" in
    remove)
        echo "🗑️ Package removal - cleaning up services and configurations..."
        
        # Stop and disable services
        systemctl stop threepics-backend.service 2>/dev/null || true
        systemctl stop threepics-frontend.service 2>/dev/null || true
        systemctl disable threepics-backend.service 2>/dev/null || true
        systemctl disable threepics-frontend.service 2>/dev/null || true
        systemctl daemon-reload
        
        # Remove cronjobs for threepics user
        if id -u threepics &>/dev/null; then
            echo "📆 Removing cronjobs for user threepics..."
            crontab -u threepics -r 2>/dev/null || true
        fi
        
        # Reset plymouth theme to default
        echo "🖼️ Resetting Plymouth theme to default..."
        plymouth-set-default-theme -R spinner 2>/dev/null || true
        
        # Reset boot config changes
        CONFIG_FILE="/boot/firmware/config.txt"
        if [ -f "$CONFIG_FILE" ]; then
            echo "🔧 Removing disable_splash setting from $CONFIG_FILE..."
            sed -i '/^disable_splash=1$/d' "$CONFIG_FILE" || true
        fi
        
        echo "✅ Package removal cleanup completed"
        ;;
        
    purge)
        echo "🧹 Package purge - removing all data and user account..."
        
        # Stop and disable services (if still running)
        systemctl stop threepics-backend.service 2>/dev/null || true
        systemctl stop threepics-frontend.service 2>/dev/null || true
        systemctl disable threepics-backend.service 2>/dev/null || true
        systemctl disable threepics-frontend.service 2>/dev/null || true
        systemctl daemon-reload
        
        # Remove threepics user and home directory
        if id -u threepics &>/dev/null; then
            echo "👤 Removing user 'threepics' and home directory..."
            userdel -r threepics 2>/dev/null || true
        fi
        
        # Remove application directory
        if [ -d "/opt/threepics" ]; then
            echo "📁 Removing /opt/threepics directory..."
            rm -rf /opt/threepics
        fi
        
        # Remove log files
        rm -f /var/log/threepics-*.log
        
        # Reset plymouth theme
        plymouth-set-default-theme -R spinner 2>/dev/null || true
        
        # Reset boot config
        CONFIG_FILE="/boot/firmware/config.txt"
        if [ -f "$CONFIG_FILE" ]; then
            sed -i '/^disable_splash=1$/d' "$CONFIG_FILE" || true
        fi
        
        echo "✅ Package purge completed - all threepics data removed"
        ;;
        
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        echo "🔄 Upgrade/abort operation detected - preserving all data"
        echo "   Argument: $1"
        echo "   Keeping /opt/threepics, user account, and configurations intact"
        echo "   Services will be managed by new package version"
        # DO NOTHING during upgrades - let the new package handle everything
        ;;
        
    *)
        echo "⚠️ postrm called with unknown argument: $1"
        # Don't exit with error for unknown arguments to be future-compatible
        ;;
esac

echo "🎯 postrm script completed at $(date)"
exit 0
